# config/wrt/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown

# comment by dimmalex for he console
#::askconsole:/bin/hetui
#::askconsole:/usr/libexec/login.sh
