# config/host/rc.local
#!/bin/bash

# add by dimmalex
runuser -u ubuntu /usr/prj/setup.sh

